# templates/default/logstash-upstart.conf.erb
#
# Generated by Chef
#
# logstash - Logstash Service
#
description     "Logstash"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

# Increase default timeout before killing (on stop) to 10 seconds
kill timeout 50

setuid <%= node['logstash']['user'] %>
setgid <%= node['logstash']['group'] %>
respawn limit 15 5
kill signal HUP
umask 0077

pre-start script
  test -f <%= @logstash_jar %> || { stop; exit 0; }
  test -d <%= node['logstash']['config_path'] %> || { stop; exit 0; }
  test -d <%= node['logstash']['log_path'] %> || { stop; exit 0; }
end script

<%
   args = []
   args << '-server'
   args << "-Xms#{node['logstash']['xms']}"
   args << "-Xmx#{node['logstash']['xmx']}"
   args << node['logstash']['java_opts'] if node['logstash']['java_opts']
   args << node['logstash']['gc_opts'] if node['logstash']['gc_opts']
   args << '-Djava.net.preferIPv4Stack=true' if node['logstash']['authbind_enable']
   args << "-Djava.io.tmpdir=#{node['logstash']['tmp_path']}"
   args << '-jar'
   args << @logstash_jar
   args << 'agent'
   args << '--config'
   args << node['logstash']['config_path']
   args << '--log'
   args << "#{node['logstash']['log_path']}/logstash.log"
%>

exec <%= node['logstash']['authbind_enable'] ? '/usr/bin/authbind --deep ' : '' %><%= node['java']['java_home'] %>/bin/java <%= args.join(' ') %>
